# challenges/core/templates/sshd/onbuild/overlay/etc/update-motd.d/50-challenge
#!/bin/sh

echo
if [ -x ~bobby/level ]; then
    printf '/--------------------------------------------------------\\\n'
    printf '|     Challenge                                          |\n'
    printf '\\--------------------------------------------------------/\n'
    echo
    ( cd ~bobby; file ./level )
    echo
fi
