# level/web/testing/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:ttyjs]
command=/usr/local/bin/tty.js --config /srv/shelljs/config.json --port 80
